# scripts/systemd/ssh-ca.service
[Unit]
Description=SSH CA Server
Documentation=https://github.com/adamscao/caserver
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Type=simple
User=ssh-ca
Group=ssh-ca

# Working directory
WorkingDirectory=/var/lib/ssh-ca

# Start the service
ExecStart=/usr/local/bin/caserver -config /etc/ssh-ca/config.yaml

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitInterval=600s
StartLimitBurst=5

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

# Read-write access to necessary paths
ReadWritePaths=/var/lib/ssh-ca /etc/ssl/ssh-ca /var/log/ssh-ca

# System call filtering (whitelist approach)
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

# Capability restrictions
CapabilityBoundingSet=
AmbientCapabilities=

# Network restrictions (allow outbound for potential future features)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssh-ca

# Watchdog (optional - requires support in application)
# WatchdogSec=30s

[Install]
WantedBy=multi-user.target
